# Agent/deltaagent.service
[Unit]
Description=deltaagent
After=syslog.target

[Service]
ExecStart=/var/datavail/agent.jar

[Install]
WantedBy=multi-user.target