# Agent/deltaagentupdater.service
[Unit]
Description=deltaagentupdater
After=syslog.target

[Service]
ExecStart=/var/datavail/agentupdater.jar

[Install]
WantedBy=multi-user.target